# cloudflared
#!/sbin/openrc-run

name="cloudflared"
command="/etc/cloudflared/cloudflared"
command_args="--config /etc/cloudflared/config.yml &"

depend() {
	need net localmount sshd
}
